# crontab
0 0 * * * python /src/manage.py flushexpiredtokens